# exam/exam-app.service
[Unit]
Description=Exam Streamlit Application (試験問題作成・採点・評価システム)
After=network-online.target
Wants=network-online.target

[Service]
Type=simple
WorkingDirectory=/opt/exam
Environment=PYTHONUNBUFFERED=1
EnvironmentFile=-/opt/exam/.env
ExecStart=/opt/exam/venv/bin/streamlit run /opt/exam/exam_app.py \
    --server.port=8505 \
    --server.address=0.0.0.0 \
    --server.headless=true
Restart=on-failure
RestartSec=5
TimeoutStopSec=30

[Install]
WantedBy=multi-user.target
