# deploy/systemd/gemini-ui2-admin.service
[Unit]
Description=Gemini Chat Admin Panel v2 (Port 8513)
After=network.target

[Service]
User=root
WorkingDirectory=/opt/gemini-ui2
Environment="PYTHONUNBUFFERED=1"
ExecStart=/opt/gemini-ui2/venv/bin/streamlit run admin.py \
    --server.port 8513 \
    --server.address 172.16.16.10 \
    --server.headless true \
    --browser.gatherUsageStats false
Restart=always
RestartSec=5

[Install]
WantedBy=multi-user.target
