# deploy/systemd/gemini-ui.service
[Unit]
Description=Gemini Chat UI (Port 8507)
After=network.target

[Service]
User=root
WorkingDirectory=/opt/gemini-ui
Environment="PYTHONUNBUFFERED=1"
ExecStart=/opt/gemini-ui/venv/bin/streamlit run app.py \
    --server.port 8507 \
    --server.address 172.16.16.10 \
    --server.headless true \
    --browser.gatherUsageStats false \
    --server.maxUploadSize 500
Restart=always
RestartSec=5

[Install]
WantedBy=multi-user.target
