# deploy/systemd/mescheck.service
[Unit]
Description=MesCheck - Microsoft 365 Message Viewer (Port 8511)
After=network-online.target
Wants=network-online.target

[Service]
Type=simple
WorkingDirectory=/opt/mescheck
Environment=PYTHONUNBUFFERED=1
EnvironmentFile=-/opt/mescheck/.env
ExecStart=/opt/mescheck/venv/bin/streamlit run /opt/mescheck/app.py \
    --server.port=8511 \
    --server.address=0.0.0.0 \
    --server.headless=true
Restart=on-failure
RestartSec=5
TimeoutStopSec=30

[Install]
WantedBy=multi-user.target
