# deploy/systemd/tts.service
[Unit]
Description=Voice Synthesis Tool (Port 8509)
After=network.target

[Service]
User=root
WorkingDirectory=/opt/tts
Environment="PYTHONUNBUFFERED=1"
ExecStart=/opt/tts/venv/bin/streamlit run app.py \
    --server.port 8509 \
    --server.address 172.16.16.10 \
    --server.headless true \
    --browser.gatherUsageStats false
Restart=always
RestartSec=5

[Install]
WantedBy=multi-user.target
